# launchd/com.oakfootwear.reporte.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oakfootwear.reporte</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/luisdarool12/Documents/OAK FOOTWEAR /WORK FLOWS /WF CONCILIACION DE CUENTAS /facturas-control/reporte_diario.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/luisdarool12/Documents/OAK FOOTWEAR /WORK FLOWS /WF CONCILIACION DE CUENTAS /facturas-control/data/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/luisdarool12/Documents/OAK FOOTWEAR /WORK FLOWS /WF CONCILIACION DE CUENTAS /facturas-control/data/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/luisdarool12/Library/Python/3.9/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
